# utils/systemd/rymon.service
[Unit]
Description=RYMonv3 dashboard backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/RYMonv3
ExecStart=/opt/RYMonv3/.venv/bin/python /opt/RYMonv3/monitor.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target